# devel_isolated/rosbash/etc/catkin/profile.d/15.rosbash.zsh
# generated from rosbash/env-hooks/15.rosbash.zsh.em

. "/opt/ros_catkin_ws/src/ros/rosbash/roszsh"
